# modules/system/host/beacon-boot/mkosi.extra/usr/local/lib/systemd/system/beacon-boot.service
[Unit]
Description=Run beacon boot in host
After=boot.target
Requires=boot.target

[Service]
Type=oneshot
ExecStart=beacon boot
StandardOutput=journal
StandardError=journal
