# modules/build/host/targets/mkosi.extra/usr/local/lib/systemd/system/stop.target
# Target reached when all the tests are complete.
# Job services should use: WantedBy=test.target

[Unit]
Description=test stage

Wants=stop-done.service
After=stop-done.service

After=report.target
Requires=report.target

[Install]
WantedBy=multi-user.target
