# modules/build/guest/targets/mkosi.extra/usr/local/lib/systemd/system/test.target
# Target reached when all the tests are complete.
# Job services should use: WantedBy=test.target

[Unit]
Description=test stage

Wants=test-done.service
After=test-done.service

After=system.target
Requires=system.target
